# board/tezuka/common/overlay_base/etc/profile.d/tailscale.sh
# shellcheck shell=sh
# Tailscale login status — disabled for now
